# deploy/harvester/polymarket-harvester.service
[Unit]
Description=Polymarket HFT Data Harvester
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/polymarket
EnvironmentFile=-/etc/default/polymarket-harvester
Environment=POLYMARKET_DATA_DIR=/opt/polymarket/data
ExecStart=/opt/polymarket/build/src/harvester/polymarket_harvester
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
