# deploy/coinbase/coinbase-harvester.service
[Unit]
Description=Coinbase BTC Price Harvester
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/polymarket
EnvironmentFile=-/etc/default/coinbase-harvester
Environment=COINBASE_DATA_DIR=/opt/polymarket/data
ExecStart=/opt/polymarket/build/src/coinbase/coinbase_harvester
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
